# systemd/repo-sync-summary.timer
[Unit]
Description=Run repo-sync weekly summary every Monday at 12:00

[Timer]
OnCalendar=*-*-* 13:00:00
Persistent=true

[Install]
WantedBy=timers.target
